# scripts/services/lelamp-setup-ap.service
[Unit]
Description=LeLamp Setup AP Mode
Documentation=https://github.com/xwang152-jack/lelamp_runtime
After=network.target
Before=lelamp-captive-portal.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/lelamp-start-ap
ExecStop=/usr/local/bin/lelamp-stop-ap
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lelamp-setup-ap

[Install]
WantedBy=multi-user.target
